# infra/systemd/conversease-web.service
[Unit]
Description=Conversease Web
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/conversease/app/apps/web
EnvironmentFile=/var/www/conversease/app/.env.production
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/bin/env npm run start -- --port 3010 --hostname 127.0.0.1
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
